# yify-downloader.service
[Unit]
Description=YIFY Torrent Downloader Service
After=network.target

[Service]
User=pi
Group=pi
WorkingDirectory=/opt/yify_downloader
ExecStart=/opt/yify_downloader/venv/bin/python /opt/yify_downloader/downloader.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Environment=DOWNLOAD_PATH=/mnt/media/downloads/movies

# Basic security measures
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target